# plugins/devops-tools/skills/pushover-verbatim-notify/templates/com.terryli.pushover-quota.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily Pushover quota monitor (iter 12b, 2026-05-19).

  Runs every day at 03:30 (15 min before pushover-prune). Calls Pushover's
  /1/apps/limits.json, persists JSON to ~/.local/state/pushover/quota.json,
  and fires a WARN-level Pushover alert when remaining < 20% of monthly cap.

  Install:
    cp templates/com.terryli.pushover-quota.plist ~/Library/LaunchAgents/
    mkdir -p ~/.local/state/launchd-logs/pushover-quota
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.terryli.pushover-quota.plist

  Verify:
    launchctl kickstart -p gui/$(id -u)/com.terryli.pushover-quota
    tail -1 ~/.local/state/launchd-logs/pushover-quota/stdout.log
    cat ~/.local/state/pushover/quota.json

  Uninstall:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.terryli.pushover-quota.plist
    rm ~/Library/LaunchAgents/com.terryli.pushover-quota.plist

  Tune alert threshold: edit --alert-threshold (0.20 = alert at 20% remaining).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.terryli.pushover-quota</string>
    <key>Program</key>
    <string>/Users/terryli/.local/bin/pushover-quota</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/terryli/.local/bin/pushover-quota</string>
        <string>--alert-threshold</string>
        <string>0.20</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/terryli/.local/state/launchd-logs/pushover-quota/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/terryli/.local/state/launchd-logs/pushover-quota/stderr.log</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/terryli</string>
        <key>PATH</key>
        <string>/Users/terryli/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
